# templates/guests/coreos/etcd.service.erb
[Unit]
Description=Clustered etcd
#After=docker.service

[Service]
Restart=always
<<<<<<< HEAD
ExecStart=/usr/bin/etcd -c 4001 -s 7001 -h <%= options[:my_ip] %> <% if options[:connection_string] %>-C <%= options[:connection_string] %><% end %> -d /home/core/etcd
=======
ExecStart=/usr/bin/etcd -c 4001 -s 7001 <% if options[:connect_to] %>-C <%= options[:connect_to] %>:7001<% end %> -d /home/core/etcd
>>>>>>> LVP for coreos guest support

[Install]
WantedBy=local.target
